# deploy/systemd/homewizard-dashboard.service
[Unit]
Description=HomeWizard web dashboard
After=network-online.target homewizard-collector.service
Wants=network-online.target

[Service]
Type=simple
User=homewizard
Group=homewizard
WorkingDirectory=/opt/homewizard
Environment=PYTHONUNBUFFERED=1
Environment=PORT=5001
ExecStart=/opt/homewizard/start_web_ui.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
